# etc/supervisord.d/kafka-offset-monitor_update.ini
[program:kafka-offset-monitor_update]
command=/opt/qnib/kafka-offset-monitor/bin/start_monitor_update.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
