# etc/supervisord.d/kafka-offset-monitor_update.ini
[program:kafka-offset-monitor_update]
command=consul-template -consul localhost:8500 -template "/etc/consul-templates/start_monitor.sh.ctmpl:/opt/qnib/kafka-offset-monitor/bin/start_monitor.sh: chmod +x opt/qnib/kafka-offset-monitor/bin/start_monitor.sh; supervisorctl restart kafka-monitor"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
